# series_notify.cron
* * * * * DISPLAY=:0.0 /home/webwerks/series_notify.py
